# tutorial/scalbf-scidb/conf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/opt/supervisord.log

[program:containerstartup]
command=/opt/container_startup.sh

